# tools/linear-import/deploy/kanon-shadow-refresh.timer
# kanon-shadow-refresh.timer — periodic shadow re-import.
# Cadence pairs with the shadow server's KANON_RELOAD_INTERVAL: the timer pulls
# from Linear on each OnCalendar tick; the server reflects appended events
# within one reload interval. A realtime (OnCalendar) schedule is used so
# Persistent=true is meaningful — a run missed while the box was down fires once
# at the next boot (monotonic OnUnitActiveSec timers get no such catch-up).
[Unit]
Description=Periodic Kanon shadow-mirror refresh (Linear -> Kanon)

[Timer]
OnCalendar=*:0/30
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
